# Linux/Week_5_Homework_Submission/cron/crontabs/root
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontabcj3xSz installed on Sat Jan 25 19:04:42 2020)
# (Cron version -- $Id: crontab.c,v 2.13 1994/01/17 03:20:37 vixie Exp $)
#Ansible: check for rootkits
@daily bash /opt/chkrootkit
#Ansible: Backup
*/2 * * * * /usr/local/bin/trigger_backup
