# Linux/Week_5_Homework_Submission/cron/crontabs/sysadmin
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.BXJNmo/crontab installed on Wed Mar 18 13:26:10 2020)
# (Cron version -- $Id: crontab.c,v 2.13 1994/01/17 03:20:37 vixie Exp $)
# Edit this file to introduce tasks to be run by cron.
# 
# Each task to run has to be defined through a single line
# indicating with different fields when the task will be run
# and what command to run for the task
# 
# To define the time you can provide concrete values for
# minute (m), hour (h), day of month (dom), month (mon),
# and day of week (dow) or use '*' in these fields (for 'any').# 
# Notice that tasks will be started based on the cron's system
# daemon's notion of time and timezones.
# 
# Output of the crontab jobs (including errors) is sent through
# email to the user the crontab file belongs to (unless redirected).
# 
# For example, you can run a backup of all your user accounts
# at 5 a.m every week with:
# 0 5 * * 1 tar -zcf /var/backups/home.tgz /home/
# 
# For more information see the manual pages of crontab(5) and cron(8)
# 
# m h  dom mon dow   command
0 18 * * * mv /home/sysadmin/Downloads/doctors*.docx /usr/share/doctors
0 18 * * * mv /home/sysadmin/Downloads/patients*.docx /usr/share/patients
0 18 * * * mv /home/sysadmin/Downloads/treatments.txt /usr/share/treatments
0 23 * * 5 tar cvf /home/sysadmin/Documents/MedicalArchive/Medicalresearch
5 23 * * 5 gzip -t /home/sysadmin/Documents/MedicalArchive/Medical_backup.gz > /home/sysadmin/backup_validation.txt
0 4 * * * ls -lR /home/sysadmin/Downloads > /home/sysadmin/Documents/Medical_files_list.txt
0 6 * * 3 tar czvf /auth_Backup.tgz /var/log/auth.log
